# overlay/etc/rc.local
#!/bin/sh

# for LicheeConsole4A
if (cat /proc/device-tree/model | grep LicheeConsole4A);
then
	echo 1 > /sys/class/graphics/fbcon/rotate_all
fi

# bluetooth driver
/usr/bin/hciattach ttyS4 any 1500000

# mipi2edp
modprobe lt8911exb

exit 0
